# tag-linux/config/autostart/xcape.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=xcape
Comment=Tap control for escape
Exec=/usr/bin/xcape
StartupNotify=false
Terminal=false
